# sillytavern/content/service/cloudflared/run
#!/bin/sh

exec cloudflared tunnel --no-autoupdate run